# root/etc/s6-overlay/s6-rc.d/postgres/run
#!/command/execlineb -P
/usr/bin/postgres
