# root/etc/s6-overlay/s6-rc.d/watcher/run
#!/command/execlineb -P

/envelope/bin/scan
