# root/etc/s6-overlay/s6-rc.d/mountebank/run
#!/command/execlineb -P
cd /work
/usr/local/bin/mb start \
    --allowInjection \
    --localOnly \
    --protofile /etc/mountebank/protocols.json \
    --configfile /etc/mountebank/all.config.json
