# root/etc/s6-overlay/s6-rc.d/mountebank/run
#!/command/execlineb -P
/envelope/bin/run-mountebank
